# supervisor.d/ape_server.ini
[program:ape-server]
command=/ape-cage/ape_server.sh
exitcodes=0
autostart=true
stdout_logfile=/ape-cage/var/apes_server.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=2
stdout_capture_maxbytes=1MB
stderr_logfile=/ape-cage/var/apes_server.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=2
stderr_capture_maxbytes=1MB
