# supervisor.d/web_app.ini
[program:web-app]
command=/usr/bin/java -jar /sample-app.jar
exitcodes=0
autostart=true
stdout_logfile=/ape-cage/var/web-app.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=2
stdout_capture_maxbytes=1MB
stderr_logfile=/ape-cage/var/web-app.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=2
stderr_capture_maxbytes=1MB
